# files/usr/etc/hotplug.d/iface/80-logger
#!/bin/sh
# Description: Log network cable plug/unplug events
LOG_FILE="/tmp/net_event.log"
echo "$(date) - ACTION: $ACTION, INTERFACE: $INTERFACE" >> $LOG_FILE